# systemd/birdcam-cleanup.timer
[Unit]
Description=Run Birdcam cleanup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
