# systemd/birdcam-stream.service
[Unit]
Description=Birdcam camera streaming pipeline
After=network.target

[Service]
Type=simple
User=birdcam
Group=birdcam
ExecStartPre=/bin/bash -c 'pkill -9 -f "rpicam-vid.*--camera 0" 2>/dev/null; sleep 1; exit 0'
ExecStart=/opt/birdcam/src/stream.sh
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=on-failure
RestartSec=3
StartLimitBurst=5
StartLimitIntervalSec=300
StandardOutput=append:/var/log/birdcam/stream.log
StandardError=append:/var/log/birdcam/stream.log

[Install]
WantedBy=multi-user.target
